# systemd/steamtrain.service
[Unit]
Description=Set hardware-appropriate Steam launch options for installed games
Documentation=https://github.com/twigglits/steam-launch-options-bot

[Service]
Type=oneshot
ExecStart=%h/.local/bin/steamtrain apply
Nice=10
# Writing is safely skipped while Steam runs; the timer retries later.
SuccessExitStatus=0
